# 200-usb-dmenu.rules
ACTION=="add", SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", RUN+="/bin/sh /home/andre/dev/usb-dmenu/usb-mgr %E{DEVNAME} %E{ID_VENDOR}"
